# amule-tools/src/proton-amule-port-sync.service
[Unit]
Description=Keep ProtonVPN forwarded port synced into aMule config
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/proton-amule-port-sync.env
ExecStart=/usr/local/bin/proton-amule-port-sync.sh
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/run /home/pi/.aMule
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
